# packages/dev-db/couchdb/files/systemd/couchdb.service
[Unit]
Description=CouchDB database
After=network.target

[Service]
User=couchdb
PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir /run/couchdb
ExecStartPre=/usr/bin/chown couchdb:couchdb /run/couchdb
ExecStart=/usr/bin/couchdb

[Install]
WantedBy=multi-user.target
