# packages/dev-db/couchdb/files/systemd/couchdb.service
[Unit]
Description=CouchDB database
After=network.target

[Service]
User=couchdb
ExecStart=/usr/bin/couchdb

[Install]
WantedBy=multi-user.target
