# packages/net-misc/dropbear/files/systemd/dropbear.service
[Unit]
Description=Dropbear SSH daemon
After=network.target

[Service]
EnvironmentFile=-/etc/default/dropbear
ExecStart=/usr/sbin/dropbear $DROPBEAR_OPTIONS -F

[Install]
WantedBy=multi-user.target
