# packages/net-misc/radvd/files/systemd/radvd.service
[Unit]
Description=IPv6 Router Advertisement Daemon
Wants=network.target
After=syslog.target

[Service]
ExecStart=/usr/sbin/radvd --nodaemon --username radvd

[Install]
WantedBy=multi-user.target
