# packages/www-servers/pound/files/systemd/pound.service
[Unit]
Description=Pound daemon
After=network.target

[Service]
Type=forking
PIDFile=/run/pound.pid
ExecStart=/usr/sbin/pound -f /etc/pound.cfg -p /run/pound.pid
KillSignal=SIGHUP

[Install]
WantedBy=multi-user.target
